# aws-rails-standalone/roles/ruby/templates/profile.d/rails_key.sh
export RAILS_ENV={{ rails_env }}
export SECRET_KEY_BASE={{ rails_secret_key_base }}
